# Systemd.service
# systemd.service - Systemd service for non-Docker deployment

[Unit]
Description=Trading Leverage API
After=network.target
Wants=network.target

[Service]
Type=exec
User=trading
Group=trading
WorkingDirectory=/opt/trading-api
Environment=FLASK_ENV=production
Environment=DATABASE_PATH=/opt/trading-api/data/trading_data.db
Environment=LOG_FILE=/opt/trading-api/logs/trading_api.log
ExecStart=/opt/trading-api/venv/bin/gunicorn --bind 0.0.0.0:8079 --workers 4 trading_api:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-api

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/trading-api/data /opt/trading-api/logs
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target